# supervisor_config.ini
[program:shadowsocks]
command=ssserver -c /root/config/ss.json  ; 需要执行的命令 wd
user =root  ;
autostart=true  ;
autorestart=true  ;
startsecs=3  ;)
stderr_logfile=/root/config/supervisor_config/logs/shadowsocks_err.log ; 错误输出重定向
stdout_logfile=None  ;log 输出
