# supervisor_config/aria2.ini
[program:aria2]
command=aria2c --conf-path="/root/config/aria2.conf"  ;  wd
user =root  ;
autostart=true  ;
autorestart=true  ;
startsecs=3  ;)
stderr_logfile=/root/config/supervisor_config/logs/ariac2_err.log ;
stdout_logfile=None  ;log
