# supervisor_config/shadowsocks.ini
[program:shadowsocks]
command=ssserver -c /root/config/ss.json  ;  wd
user =root  ;
autostart=true  ;
autorestart=true  ;
startsecs=3  ;)
stderr_logfile=/root/config/supervisor_config/logs/shadowsocks_err.log ;
stdout_logfile=None  ;log
